# lib/cxl_cli/cxl/cxl-monitor.service
[Unit]
Description=CXL Monitor Daemon

[Service]
Type=simple
ExecStart=/usr/bin/cxl monitor

[Install]
WantedBy=multi-user.target
