# lib/cxl_cli/ndctl/ndctl-monitor.service
[Unit]
Description=Ndctl Monitor Daemon

[Service]
Type=simple
ExecStart=/usr/bin/ndctl monitor

[Install]
WantedBy=multi-user.target
